# bookmark-bot-update.timer
[Unit]
Description=Run update-check on bookmark-botd every 2 hour

[Timer]
OnBootSec=15min
OnUnitActiveSec=2hour

[Install]
WantedBy=timers.target
